# ansible/roles/switchbot_bridge/templates/switchbot-bridge.service.j2
[Unit]
Description=SwitchBot BLE to MQTT bridge
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
User=switchbot-bridge
Group=switchbot-bridge
ExecStart=/opt/switchbot-bridge/venv/bin/python /opt/switchbot-bridge/switchbot_bridge.py
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false

[Install]
WantedBy=multi-user.target
