# ansible/roles/hermes_agent/templates/hermes-gateway.service.j2
[Unit]
Description=Hermes Agent Gateway
After=network.target docker.service nftables.service
Requires=docker.service

[Service]
Type=simple
User={{ hermes_user }}
Group={{ hermes_group }}
WorkingDirectory={{ hermes_home }}
Environment=HERMES_HOME={{ hermes_home }}
Environment=HOME={{ hermes_home }}
EnvironmentFile={{ hermes_config_dir }}/secrets.env
ExecStart={{ hermes_binary_link }} gateway start --config {{ hermes_home }}/config.yaml
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
