# ansible/roles/glow_mqtt_bridge/templates/glow-mqtt-bridge.service.j2
[Unit]
Description=Glowmarkt (UK smart meter, DCC) to MQTT bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=glow-mqtt-bridge
Group=glow-mqtt-bridge
ExecStart=/opt/glow-mqtt-bridge/venv/bin/python /opt/glow-mqtt-bridge/glow_mqtt_bridge.py
Restart=on-failure
RestartSec=5

# systemd creates /var/lib/glow-mqtt-bridge (owned by the service user) and
# exports STATE_DIRECTORY; the daemon caches its auth token there.
StateDirectory=glow-mqtt-bridge

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false

[Install]
WantedBy=multi-user.target
